# deploy/trojan-auto-update.service
[Unit]
Description=Refresh Trojan Docker stack from GitHub
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trajondocker
ExecStart=/opt/trajondocker/scripts/update.sh
